# scripts/unified/out/story-goal-mcp/linux-systemd/story-goal-mcp.service
[Unit]
Description=story-goal-mcp (MCP stdio)
After=network.target

[Service]
Type=simple
WorkingDirectory=/Users/hallie/Documents/repos/tools/lisa_brain
Environment=PYTHONUNBUFFERED=1
ExecStart=/usr/bin/env bash -lc 'python3 /Users/hallie/Documents/repos/tools/lisa_brain/story_goal_mcp.py'
Restart=on-failure
RestartSec=2
User=www-data
Group=www-data

[Install]
WantedBy=multi-user.target
